# root/etc/services.d/d-zone/run
#!/usr/bin/with-contenv bash

cd /app/d-zone || exit

nginx

exec \
	s6-setuidgid abc /usr/bin/npm start
